# board/raspberrypi/overlay/usr/lib/systemd/system/rfkill_unblock.service
[Unit]
Description=Disable rfkill block

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/rfkill unblock all
