# board/raspberrypi/overlay/usr/lib/systemd/system/wpa_supplicant.service
[Unit]
Description=WPA supplicant

[Service]
Type=exec
ExecStart=/sbin/wpa_supplicant -i wlan0 -c /etc/wpa_supplicant_moosephone.conf
