# board/raspberrypi/overlay/lib/systemd/system/rfkill_unblock.service
[Unit]
Description=Disable rfkill block

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/rfkill unblock all
